# deploy/systemd/helmsman.service
# Helmsman runs as its own systemd unit, NOT a compose container, so it can
# never appear in a managed project's container list and a stack `down` can't
# take it down (plan §3). Runs non-root as a dedicated low-priv user in the
# `docker` group. Hardening per plan §3 / §15 Phase 4.
#
# Install:
#   useradd --system --no-create-home --shell /usr/sbin/nologin helmsman
#   usermod -aG docker helmsman
#   install -m0755 helmsman /usr/local/bin/helmsman
#   install -d -o helmsman -g helmsman -m0700 /var/lib/helmsman
#   install -d -o helmsman -g helmsman -m0700 /var/lib/helmsman-apps  # app run dirs
#   install -d -o helmsman -g helmsman -m0700 /var/lib/caddy          # edge data/cert store
#   # Config is root-OWNED and group-readable by the service group so the
#   # non-root service can read it; the master key is never group-WRITABLE:
#   install -o root -g helmsman -m0640 config.yaml /etc/helmsman/config.yaml
#   systemctl daemon-reload && systemctl enable --now helmsman

[Unit]
Description=Helmsman ops dashboard (control plane)
Documentation=https://github.com/daboss2003/Helmsman
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=helmsman
Group=helmsman
SupplementaryGroups=docker
ExecStart=/usr/local/bin/helmsman serve --config /etc/helmsman/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2s

# --- Memory accounting: MemoryMax kills INSIDE this cgroup, incl. the supervised
# Caddy + nginx children AND forked docker/compose (plan §0, §3). They are CO-RESIDENT
# in this cgroup (not a separate slice), so the cap must cover all of them: 192M with a
# 150M Go limit left almost nothing for Caddy+nginx and risked OOM-killing the control
# plane itself. GOMEMLIMIT sits under the cap for the Go process.
MemoryMax=384M
Environment=GOMEMLIMIT=150MiB
# OOMScoreAdjust only biases the GLOBAL killer; it is NOT a substitute for MemoryMax.
OOMScoreAdjust=-100
# HOME for forked children only (Helmsman itself never reads HOME): docker compose's
# buildx/BuildKit writes its cache to $HOME/.cache and its config to $HOME/.docker, and
# git writes $HOME/.gitconfig — all must land in a writable dir or `--build` deploys fail.
Environment=HOME=/var/lib/helmsman

# --- Sandbox (plan §15 Phase 4). The raw docker.sock is NEVER mounted here;
# reads go through the read-only socket-proxy, writes shell out to `docker compose`.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
# NOTE: @resources is intentionally NOT denied. The supervised edge/L4 children
# (Caddy, nginx) raise their file-descriptor limit at startup (setrlimit/prlimit64,
# which live in @resources) and inherit this filter via NoNewPrivileges — denying it
# kills them with SIGSYS ("bad system call"). @privileged/@obsolete stay denied.
SystemCallFilter=~@privileged @obsolete
# A blocked syscall returns EPERM instead of killing the process, so a child that
# probes a denied call degrades gracefully rather than core-dumping the edge.
SystemCallErrorNumber=EPERM

# Only these state dirs are writable (ProtectSystem=strict makes everything else
# read-only). The admin UI binds loopback only; Helmsman never binds a public port
# itself (the managed edge child does). The three paths, for the DEFAULT data_dir
# (/var/lib/helmsman) — if you change data_dir, update the first two to match:
#   /var/lib/helmsman        DB, keys, git object store
#   /var/lib/helmsman-apps   per-app run dirs (a sibling of data_dir by design — see
#                            appsRoot: keeps app binds from reaching the DB/keys)
#   /var/lib/caddy           the supervised edge's data/cert store (HOME/XDG_DATA_HOME)
# Omit any of them and the corresponding writes fail (Caddy/nginx exit, deploys can't
# materialize) — the failure mode that looks like "edge won't start / deploy hangs".
ReadWritePaths=/var/lib/helmsman /var/lib/helmsman-apps /var/lib/caddy

# Runtime dir: systemd creates /run/helmsman (owned helmsman:helmsman, writable even
# under ProtectSystem=strict, torn down on stop) BEFORE ExecStart. The managed edge's
# Caddy admin socket lives at unix//run/helmsman/caddy-admin.sock — without this the
# child can't bind it and the edge crash-loops. Do NOT mkdir /run/helmsman by hand
# (it's tmpfs, lost on reboot); this directive is the durable creator.
RuntimeDirectory=helmsman
RuntimeDirectoryMode=0700

# CAP_NET_BIND_SERVICE lets the managed edge/L4 children (Caddy :80/:443, nginx :53)
# bind privileged ports while staying NON-ROOT. Granted by DEFAULT so a managed edge
# works out of the box — no drop-in, no setup step. It is the ONLY allowed capability
# (bounding set), and ambient so the forked children inherit it. The children are
# co-resident in this unit, so the control-plane process holds it too (a child-only
# grant via a separate edge unit is a planned hardening — docs/architecture.md). In
# `external` mode the edge isn't started and this cap is simply unused (override with
# a drop-in `AmbientCapabilities=` / `CapabilityBoundingSet=` if you want zero caps).
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# --- Egress allow-listing (OPT-IN hardening, OFF by default). The primary SSRF/
# call-home defense is in-process: pinned dialers reject loopback/link-local/metadata
# for the ops prober, edge upstreams, and alert notifiers. A systemd cgroup egress
# filter is a deeper backstop, but it is OFF by default because it easily breaks a
# managed edge: Let's Encrypt has NO fixed CIDR (CDN-fronted), so IPAddressDeny=any
# would block ACME cert issuance, and localhost-only also blocks proxying to app
# containers on the docker bridge. To enable it, uncomment + tune to YOUR host:
#   IPAddressDeny=any
#   IPAddressAllow=localhost
#   IPAddressAllow=172.16.0.0/12        # docker bridge / app subnet (narrow if you can)
#   IPAddressAllow=<resolver-ip>        # only if /etc/resolv.conf is NOT 127.0.0.53
#   # ...plus an ACME-reachable egress path (impractical to pin to a CIDR — review
#   #    your CA/CDN, or keep this OFF and rely on the in-process dialers).

[Install]
WantedBy=multi-user.target
